# openrc/cloudflared
#!/sbin/openrc-run
name="cloudflared"
description="Cloudflare Tunnel"
command="/usr/bin/tmux"
command_args="new-session -d -s cloudflared 'cloudflared tunnel --config /root/.cloudflared/config.yml run >> /root/logs/cloudflared.log 2>&1'"
command_background=false

start_pre() {
  mkdir -p /root/logs
  [ -f /root/.cloudflared/config.yml ] || { eerror "Config mancante: esegui blade-setup-cloudflare"; return 1; }
  if tmux has-session -t cloudflared 2>/dev/null; then
    eend 0
    return 1
  fi
}

stop() {
  tmux kill-session -t cloudflared 2>/dev/null || true
}
