# openrc/opencode
#!/sbin/openrc-run
name="opencode"
description="OpenCode AI Server"
command="/usr/bin/tmux"
command_args="new-session -d -s opencode 'opencode serve --ui-password \"${OPENCODE_UI_PASSWORD:-}\" --hostname ${OPENCODE_HOSTNAME:-0.0.0.0} --port ${OPENCODE_PORT:-3000} >> /root/logs/opencode.log 2>&1'"
command_background=false

start_pre() {
  mkdir -p /root/logs
  if tmux has-session -t opencode 2>/dev/null; then
    eend 0
    return 1
  fi
}

stop() {
  tmux kill-session -t opencode 2>/dev/null || true
}
